# inittab
# Start "rc init" on boot
::sysinit:/etc/init.d/rc

# Set up the TTY's 1 through 4
tty1::respawn:/sbin/agetty -a me --noclear -8 -s 38400 tty1 linux
tty2::respawn:/sbin/agetty -8 -s 38400 tty2 linux
tty3::respawn:/sbin/agetty -8 -s 38400 tty3 linux
tty4::respawn:/sbin/agetty -8 -s 38400 tty4 linux

# Shutdown when pressing CTRL+ALT+DEL (disabled by default)
#::ctrlaltdel:/bin/kill -USR2 1

# Shutdown
::shutdown:/bin/killall systemd-udevd udevd 

::shutdown:echo :: running /etc/init.d/rc.shutdown if it exists
::shutdown:/etc/init.d/rc.shutdown

::shutdown:/bin/umount -r /run

::shutdown:echo :: sending SIGTERM to all
::shutdown:/bin/kill -s TERM -1
::shutdown:sleep 1
::shutdown:echo :: sending SIGKILL to all
::shutdown:/bin/kill -s KILL -1

::shutdown:echo :: unmounting everything
::shutdown:/bin/umount -a -r
::shutdown:/bin/mount -o remount,ro /
