# busybox/inittab
# boot and shutdown scripts to run, respectively
::sysinit:/etc/rc/rc.init
::shutdown:/etc/rc/rc.shutdown

# Shutdown when pressing CTRL+ALT+DEL
::ctrlaltdel:/bin/kill -USR2 1
